# deploy/systemd/btxpool-peers.timer
[Unit]
Description=Monitor BTX near-tip peers

[Timer]
OnBootSec=3min
OnUnitActiveSec=5min
Persistent=true

[Install]
WantedBy=timers.target
