# deploy/systemd/btxpool-backup.timer
[Unit]
Description=Hourly BTX pool ledger backup

[Timer]
OnCalendar=hourly
Persistent=true
RandomizedDelaySec=120

[Install]
WantedBy=timers.target
