# deploy/systemd/btxpool-peers.service
[Unit]
Description=Check BTX near-tip peer availability
After=btxd.service

[Service]
Type=oneshot
User=aravindthana
Group=aravindthana
WorkingDirectory=/mnt/e/Business/btxpool
Environment=BTXPOOL_CONFIG=/mnt/e/Business/btxpool/config.yaml
ExecStart=/mnt/e/Business/btxpool/scripts/ensure-peers.sh
