# deploy/systemd/btxpool-health.timer
[Unit]
Description=Monitor BTX pool health

[Timer]
OnBootSec=2min
OnUnitActiveSec=1min
Persistent=true

[Install]
WantedBy=timers.target
